# flare.service
# File path: /etc/systemd/system/flare.service
# sudo systemctl daemon-reload
# sudo systemctl start flare
# sudo systemctl enable flare
# sudo systemctl status flare
[Unit]
Description=Flare Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/flare  # 注意修改路径
ExecStart=/path/to/flare/flare --port 3000 --log-dir /path/to/flare/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
